# klipper_usb_autoconnect.service
[Unit]
Description=Klipper USB auto connect

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/local/bin/klipper_usb_autoconnect
